# deploy/simple-backend.service
[Unit]
Description=Simple Text Backend
After=network.target

[Service]
Type=simple
ExecStart=/opt/simple-backend/simple-backend
WorkingDirectory=/opt/simple-backend
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=simple-backend

[Install]
WantedBy=multi-user.target